# dbus/data/gpio-manager.service.in
# SPDX-License-Identifier: CC0-1.0
# SPDX-FileCopyrightText: 2023-2024 Bartosz Golaszewski <bartosz.golaszewski@linaro.org>

[Unit]
Description=Centralized GPIO manager daemon

[Service]
Type=dbus
BusName=io.gpiod1
ExecStart=@prefix@/bin/gpio-manager
Restart=always
User=gpio-manager

CapabilityBoundingSet=
ReadOnlyDirectories=/
NoNewPrivileges=yes
RemoveIPC=yes
PrivateTmp=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ProtectClock=yes
Delegate=no
IPAddressDeny=any
KeyringMode=private
LockPersonality=yes
MemoryDenyWriteExecute=yes
NotifyAccess=main
PrivateMounts=no
PrivateNetwork=no
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallFilter=~@clock
SystemCallFilter=~@cpu-emulation
SystemCallFilter=~@debug
SystemCallFilter=~@module
SystemCallFilter=~@mount
SystemCallFilter=~@obsolete
SystemCallFilter=~@privileged
SystemCallFilter=~@raw-io
SystemCallFilter=~@reboot
SystemCallFilter=~@swap

[Install]
WantedBy=multi-user.target
